# docker/s6-rc.d/nginx/run
#!/usr/bin/execlineb -P
nginx -e "stderr" -c "/run/nginx.conf" -g "daemon off;"
