# docker/s6-rc.d/blesta-init/up
#!/usr/bin/execlineb -P
elglob FILE /opt/blesta/defaults/config/* cp ${FILE} /opt/blesta/data/config/
