# docker/s6-rc.d/blesta-init/up
#!/usr/bin/execlineb -P
elglob FILE /opt/blesta/public/config-default/*.php cp ${FILE} /opt/blesta/data/config/
