# docker/s6-rc.d/vector/run
#!/usr/bin/execlineb -P
vector --quiet --config /etc/vector.toml
