# docker/s6-rc.d/tamper-guard/finish
#!/usr/bin/execlineb -P
s6-fatal "unexpected service termination: tamper-guard"
