# docker/s6-rc.d/supercronic/run
#!/usr/bin/execlineb -P
supercronic -no-reap /run/supercronic
